# deploy/com.aaitrade.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aaitrade.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/aditya/AAItrade/scripts/start_server.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/aditya/AAItrade</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/aditya/AAItrade/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/aditya/AAItrade/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/aditya/AAItrade/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
